# deploy/systemd/phone-worker-watch.timer
[Unit]
Description=Periodic phone worker helper check

[Timer]
OnBootSec=2min
OnUnitActiveSec=1min
AccuracySec=10s
Persistent=true

[Install]
WantedBy=timers.target
