# deploy/systemd/phone-worker-watch.timer
[Unit]
Description=Periodic phone worker helper check

[Timer]
OnBootSec=2min
OnUnitActiveSec=1min
AccuracySec=1s
Persistent=true
RandomizedDelaySec=0

[Install]
WantedBy=timers.target
