# deploy/systemd/tts-bot-updater.service
[Unit]
Description=Atualiza o tts-bot e reinicia o serviço se houver commit novo
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=root
Nice=10
CPUWeight=20
IOWeight=20
ExecStart=/usr/bin/env bash /home/ubuntu/bot/scripts/tts-bot-update.sh
