# deploy/systemd/vps/tts-bot.service.d/10-audio-node.conf
# Core VPS audio policy
#
# VPS-local Lavalink is intentionally not a dependency of tts-bot anymore.
# Music audio is owned by the phone worker / Music Agent side.
# Keep this drop-in in the repository so old VPS-local audio-node drop-ins are overwritten.
[Unit]
Wants=network-online.target
After=network-online.target

[Service]
# Do not wait for / start VPS-local Lavalink here.
