# deploy/systemd/vps/tts-bot.service.d/99-clear-prestart.conf
# Clear old ExecStartPre hooks left by older VPS-local Lavalink deployments.
# systemd resets list-valued ExecStartPre when it appears empty.
[Service]
ExecStartPre=
