# deploy/systemd/vps/tts-bot.service.d/zz-clear-prestart.conf
# Compatibility drop-in for VPSes that already had this file.
# It remains versioned so live /etc/systemd/system state is no longer hidden.
[Service]
ExecStartPre=
