# deploy/systemd/vps/phone-worker-watch.timer
[Unit]
Description=Periodic phone worker helper check

[Timer]
OnBootSec=2min
OnUnitActiveSec=10min
AccuracySec=1min
Persistent=true
RandomizedDelaySec=45s

[Install]
WantedBy=timers.target
