# deploy/systemd/vps/sinuca-activity-server.service
[Unit]
Description=Osaka Web Dashboard Server
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/bot
Environment=NODE_ENV=production
Environment=DASHBOARD_SERVICE_NAME=sinuca-activity-server.service
Environment=DASHBOARD_NOTIFY_COOLDOWN_SECONDS=900
ExecStartPre=/bin/bash -lc 'fuser -k 8787/tcp >/dev/null 2>&1 || true'
ExecStart=/usr/bin/env bash /home/ubuntu/bot/start-sinuca-server.sh
ExecStartPost=/usr/bin/env bash /home/ubuntu/bot/notify-activity-start.sh
Restart=on-failure
RestartSec=5
TimeoutStopSec=20
KillSignal=SIGTERM

[Install]
WantedBy=multi-user.target
